# supervisord-cruisecontrol.conf
[program:cruisecontrol]
command=/etc/init.d/cruisecontrol
numprocs=1
autostart=true
autorestart=true
